# client-vnc/configs/rootfs-v3/etc/s6-overlay/s6-rc.d/x11vnc/run
#!/command/with-contenv sh

exec /usr/bin/x11vnc -q -forever -shared -nomodtweak -capslock -noxfixes -cursor most -display ${DISPLAY}
