# swarm-jenkins-agent/configs/rootfs-v3/etc/s6-overlay/s6-rc.d/jenkins-agent/run
#!/command/with-contenv sh

wget -O agent.jar "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_JAR_URL}"
exec java -jar agent.jar \
  -url "${JENKINS_URL}" \
  -secret "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_SECRET}" \
  -noReconnect \
  -name "${JENKINS_AGENT_NAME}" \
  -webSocket \
  -workDir /home/jenkins
