# swarm-jenkins-agent/configs/rootfs-v3/etc/s6-overlay/s6-rc.d/jenkins-agent/run
#!/command/with-contenv sh

set -eu

: "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_JAR_URL:?Required env var not set}"
: "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_JNLP_URL:?Required env var not set}"
: "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_SECRET:?Required env var not set}"

wget -O agent.jar "${DOCKER_SWARM_PLUGIN_JENKINS_AGENT_JAR_URL}"
exec java -jar agent.jar \
  -jnlpUrl "${$DOCKER_SWARM_PLUGIN_JENKINS_AGENT_JNLP_URL}" \
  -secret "${$DOCKER_SWARM_PLUGIN_JENKINS_AGENT_SECRET}" \
  -noReconnect \
  -workDir /home/jenkins
